# snaps/dotnet-sdk-110/mounts/var-snap-dotnet-common-dotnet-sdk\x2dmanifests-11.0.100\x2dpreview.1-microsoft.net.sdk.android.mount
[Unit]
Description={SNAP}-sdk-manifests-microsoft.net.sdk.android

[Mount]
What=/snap/{SNAP}/current/usr/lib/dotnet/sdk-manifests/11.0.100-preview.1/microsoft.net.sdk.android
Where=/var/snap/dotnet/common/dotnet/sdk-manifests/11.0.100-preview.1/microsoft.net.sdk.android
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
