# snaps/dotnet-sdk-80/mounts/var-snap-dotnet-common-dotnet-sdk\x2dmanifests-8.0.100-microsoft.net.workload.emscripten.net6.mount
[Unit]
Description={SNAP}-sdk-manifests-microsoft.net.workload.emscripten.net6

[Mount]
What=/snap/{SNAP}/current/usr/lib/dotnet/sdk-manifests/8.0.100/microsoft.net.workload.emscripten.net6
Where=/var/snap/dotnet/common/dotnet/sdk-manifests/8.0.100/microsoft.net.workload.emscripten.net6
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
